# crontab
@daily cp -f ~/.bash_history ~/bash_history
